# roles/reboot/files/weekly_reboot.timer
[Unit]
Description=Timer for Weekly Reboot the system

[Timer]
OnCalendar=Tue *-*-* 14:25:00

