# roles/reboot/files/weekly_reboot.timer
[Unit]
Description=Timer for Weekly Reboot the system

[Timer]
OnCalendar=Sun *-*-* 00:10:00

