# roles/reboot/files/weekly_reboot.service
[Unit]
Description=Weekly System Reboot

[Service]
Type=oneshot
ExecStart=/sbin/reboot

